# deploy/systemd/rainbow-backend-test.service
[Unit]
Description=Rainbow Backend Test Service
After=network.target mysql.service
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/rainbow-backend
EnvironmentFile=/opt/rainbow-backend/test.env
ExecStart=/opt/rainbow-backend/bin/rainbow-backend
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
